# relay_hub.conf
[supervisord]
logfile = /tmp/supervisord.log
[program:relay_hub]
command=python3 /home/user/relay_hub.py
stdout_logfile=/tmp/relay_hub.log
autostart=true
autorestart=true
[supervisorctl]
;

